# transmissiond.conf
[program:transmissiond]
command = transmission-daemon -f -g /var/lib/transmission-daemon/info -v %(ENV_ADMIN_PASS)s
autorestart = true

